# backend/data/intent_metadata.json
{
  "brute_force": {
    "label": "Brute Force",
    "mitre_tactic_id": "TA0006",
    "mitre_tactic": "Credential Access",
    "mitre_technique_id": "T1110",
    "mitre_technique": "Brute Force",
    "default_severity": "High",
    "default_threshold": 5,
    "default_window_minutes": 5,
    "log_source": "Authentication logs (Windows Security 4625, Linux auth.log, Okta, Azure AD SigninLogs)",
    "required_fields": ["source_ip", "user", "timestamp", "outcome"],
    "false_positives": ["Password spraying penetration tests", "VPN clients auto-reconnecting", "Misconfigured applications retrying stale credentials", "Shared service accounts with expired passwords"],
    "coverage_gaps": ["Slow, low-and-slow brute force spread over days", "Distributed brute force from many source IPs", "Offline password attacks against exfiltrated hashes"],
    "improvements": ["Add geo-location and impossible travel enrichment", "Correlate with threat intelligence IP reputation", "Weight by asset criticality", "Incorporate a user risk score baseline"]
  },
  "credential_access": {
    "label": "Credential Access",
    "mitre_tactic_id": "TA0006",
    "mitre_tactic": "Credential Access",
    "mitre_technique_id": "T1003",
    "mitre_technique": "OS Credential Dumping",
    "default_severity": "Critical",
    "default_threshold": 1,
    "default_window_minutes": 10,
    "log_source": "Endpoint (Sysmon 10, EDR process/memory access), Windows Security logs",
    "required_fields": ["process_name", "target_process", "user", "host"],
    "false_positives": ["Legitimate backup or AV software accessing LSASS memory", "Administrative credential vaulting tools", "Sanctioned red-team exercises"],
    "coverage_gaps": ["In-memory-only dumping tools with custom binaries", "Living-off-the-land credential access techniques", "Encrypted or fileless dump exfiltration"],
    "improvements": ["Add process ancestry / parent-child correlation", "Enrich with EDR memory access alerts", "Track access to LSASS handle requests with specific access masks"]
  },
  "malware": {
    "label": "Malware",
    "mitre_tactic_id": "TA0002",
    "mitre_tactic": "Execution",
    "mitre_technique_id": "T1204",
    "mitre_technique": "User Execution",
    "default_severity": "High",
    "default_threshold": 1,
    "default_window_minutes": 5,
    "log_source": "EDR/AV alerts, Sysmon (Event ID 1, 7, 11), process creation logs",
    "required_fields": ["file_hash", "process_name", "file_path", "host"],
    "false_positives": ["Security testing tools (EICAR, custom pentest payloads)", "Legitimate but unsigned internal tooling", "Packed/compressed legitimate software"],
    "coverage_gaps": ["Zero-day or previously unseen malware without known hashes", "Fileless malware executing solely in memory", "Polymorphic malware that changes signature per infection"],
    "improvements": ["Add YARA rule based content inspection", "Enrich hashes with threat intel reputation feeds", "Add behavioral sandboxing detonation results"]
  },
  "powershell": {
    "label": "PowerShell",
    "mitre_tactic_id": "TA0002",
    "mitre_tactic": "Execution",
    "mitre_technique_id": "T1059.001",
    "mitre_technique": "Command and Scripting Interpreter: PowerShell",
    "default_severity": "Medium",
    "default_threshold": 1,
    "default_window_minutes": 5,
    "log_source": "PowerShell Script Block Logging (Event ID 4104), Sysmon Event ID 1",
    "required_fields": ["command_line", "script_block_text", "user", "host"],
    "false_positives": ["Legitimate IT automation and configuration management scripts", "Software deployment tools using PowerShell", "Developer scripting activity"],
    "coverage_gaps": ["Heavily obfuscated scripts using custom encoding beyond base64", "PowerShell run via alternate hosts (e.g. powershell.exe renamed)", "Downgrade attacks that bypass script block logging entirely"],
    "improvements": ["Decode base64 command blocks automatically for inspection", "Correlate with network connections initiated post-execution", "Baseline normal admin PowerShell usage per host"]
  },
  "lateral_movement": {
    "label": "Lateral Movement",
    "mitre_tactic_id": "TA0008",
    "mitre_tactic": "Lateral Movement",
    "mitre_technique_id": "T1021",
    "mitre_technique": "Remote Services",
    "default_severity": "High",
    "default_threshold": 3,
    "default_window_minutes": 15,
    "log_source": "Windows Security (4624 Type 3/10), Sysmon Network (3), firewall/NetFlow logs",
    "required_fields": ["source_host", "dest_host", "user", "logon_type"],
    "false_positives": ["IT administrators performing routine remote maintenance", "Automated patch management systems", "Backup agents connecting across hosts"],
    "coverage_gaps": ["Movement over uncommon or encrypted protocols", "Living-off-the-land lateral movement using native tools", "Movement within a single broadcast domain not crossing monitored chokepoints"],
    "improvements": ["Build a baseline of normal host-to-host communication graphs", "Correlate with prior reconnaissance or credential access alerts", "Add asset criticality weighting for crown-jewel systems"]
  },
  "persistence": {
    "label": "Persistence",
    "mitre_tactic_id": "TA0003",
    "mitre_tactic": "Persistence",
    "mitre_technique_id": "T1053",
    "mitre_technique": "Scheduled Task/Job",
    "default_severity": "Medium",
    "default_threshold": 1,
    "default_window_minutes": 10,
    "log_source": "Windows Security (4698, 4657), Sysmon (Event ID 13), Registry/Task Scheduler logs",
    "required_fields": ["task_name", "registry_key", "process_name", "host"],
    "false_positives": ["Legitimate software installers creating scheduled tasks", "IT-managed configuration baselines", "Update agents registering autorun entries"],
    "coverage_gaps": ["Persistence via less common autostart extensibility points (ASEPs)", "Persistence hidden in WMI repository without file artifacts", "Cloud-native persistence (e.g. serverless function triggers)"],
    "improvements": ["Maintain an allowlist of known-good scheduled tasks and services", "Correlate creation events with the parent process reputation", "Add periodic ASEP baseline diffing"]
  },
  "ransomware": {
    "label": "Ransomware",
    "mitre_tactic_id": "TA0040",
    "mitre_tactic": "Impact",
    "mitre_technique_id": "T1486",
    "mitre_technique": "Data Encrypted for Impact",
    "default_severity": "Critical",
    "default_threshold": 20,
    "default_window_minutes": 2,
    "log_source": "File integrity monitoring, Sysmon (Event ID 11, 23), EDR file activity telemetry",
    "required_fields": ["file_path", "process_name", "host", "file_extension"],
    "false_positives": ["Legitimate bulk file compression or archiving jobs", "Backup software rewriting large file sets", "File synchronization/migration tools"],
    "coverage_gaps": ["Slow encryption designed to stay under volumetric thresholds", "Ransomware that encrypts only a small set of high-value files", "Double-extortion actors who exfiltrate before encrypting, bypassing pure encryption-rate detection"],
    "improvements": ["Add shadow copy and backup deletion correlation", "Monitor file entropy changes rather than just rename volume", "Add canary/honeyfile triggers across shares"]
  },
  "phishing": {
    "label": "Phishing",
    "mitre_tactic_id": "TA0001",
    "mitre_tactic": "Initial Access",
    "mitre_technique_id": "T1566",
    "mitre_technique": "Phishing",
    "default_severity": "Medium",
    "default_threshold": 1,
    "default_window_minutes": 60,
    "log_source": "Email gateway logs, Microsoft 365/Google Workspace audit logs, proxy logs",
    "required_fields": ["sender", "recipient", "attachment_name", "url"],
    "false_positives": ["Legitimate marketing emails resembling phishing patterns", "Internal security awareness phishing simulations", "Newly registered but legitimate business domains"],
    "coverage_gaps": ["Spear-phishing with no attachment or link (pure social engineering)", "Phishing hosted on compromised legitimate domains", "Multi-stage phishing where the malicious payload is delivered later"],
    "improvements": ["Integrate with URL sandboxing/detonation services", "Add sender domain age and reputation enrichment", "Correlate with subsequent user click and authentication events"]
  },
  "privilege_escalation": {
    "label": "Privilege Escalation",
    "mitre_tactic_id": "TA0004",
    "mitre_tactic": "Privilege Escalation",
    "mitre_technique_id": "T1078",
    "mitre_technique": "Valid Accounts",
    "default_severity": "High",
    "default_threshold": 1,
    "default_window_minutes": 10,
    "log_source": "Windows Security (4672, 4732), cloud IAM audit logs (AWS CloudTrail, Azure AD Audit)",
    "required_fields": ["user", "target_group", "role", "host_or_account"],
    "false_positives": ["Approved change-management driven group membership changes", "Break-glass emergency access procedures", "Automated provisioning tools assigning roles"],
    "coverage_gaps": ["Privilege escalation via undocumented software vulnerabilities", "Escalation using legitimate but overly-permissive existing roles", "Escalation chains spanning multiple cloud accounts"],
    "improvements": ["Correlate with just-in-time access approval workflows", "Add anomaly detection for rarely-used privileged roles", "Flag privilege changes outside change windows"]
  },
  "data_exfiltration": {
    "label": "Data Exfiltration",
    "mitre_tactic_id": "TA0010",
    "mitre_tactic": "Exfiltration",
    "mitre_technique_id": "T1041",
    "mitre_technique": "Exfiltration Over C2 Channel",
    "default_severity": "Critical",
    "default_threshold": 100,
    "default_window_minutes": 10,
    "log_source": "Proxy/firewall logs, DLP alerts, cloud storage audit logs, NetFlow",
    "required_fields": ["source_host", "destination", "bytes_transferred", "user"],
    "false_positives": ["Scheduled legitimate backups to cloud storage", "Large but sanctioned data migrations", "Media-heavy business workflows (video/design teams)"],
    "coverage_gaps": ["Exfiltration disguised within normal encrypted business traffic", "Slow low-volume exfiltration over long periods", "Exfiltration via removable media without endpoint DLP coverage"],
    "improvements": ["Add destination reputation and geolocation enrichment", "Baseline per-user/per-department normal transfer volume", "Correlate with prior staging or compression activity"]
  },
  "web_attacks": {
    "label": "Web Attacks",
    "mitre_tactic_id": "TA0001",
    "mitre_tactic": "Initial Access",
    "mitre_technique_id": "T1190",
    "mitre_technique": "Exploit Public-Facing Application",
    "default_severity": "High",
    "default_threshold": 5,
    "default_window_minutes": 5,
    "log_source": "Web server/access logs, WAF logs, reverse proxy logs",
    "required_fields": ["source_ip", "uri", "http_method", "status_code"],
    "false_positives": ["Automated vulnerability scanners run by the security team", "Aggressive but legitimate web crawlers", "QA/testing traffic hitting staging endpoints"],
    "coverage_gaps": ["Zero-day exploits with no known signature", "Attacks encoded to evade WAF pattern matching", "Slow, single-request exploitation with no volumetric signal"],
    "improvements": ["Add WAF rule correlation and payload decoding", "Enrich with source IP threat intelligence", "Track sequential attack chains (recon -> exploit -> webshell)"]
  },
  "cloud_attacks": {
    "label": "Cloud Attacks",
    "mitre_tactic_id": "TA0001",
    "mitre_tactic": "Initial Access",
    "mitre_technique_id": "T1078.004",
    "mitre_technique": "Valid Accounts: Cloud Accounts",
    "default_severity": "High",
    "default_threshold": 1,
    "default_window_minutes": 15,
    "log_source": "AWS CloudTrail, Azure AD SigninLogs, GCP Audit Logs, Kubernetes audit logs",
    "required_fields": ["actor", "action", "resource", "source_ip"],
    "false_positives": ["Legitimate administrative activity from new locations", "Automation/service accounts performing bulk actions", "CI/CD pipelines provisioning cloud resources"],
    "coverage_gaps": ["Abuse of legitimate, narrowly-scoped API permissions", "Multi-cloud attack chains not correlated across providers", "Actions performed via compromised but valid service principal credentials"],
    "improvements": ["Add impossible-travel and device-trust enrichment for cloud logins", "Baseline normal API call patterns per identity", "Correlate IAM changes with subsequent resource access"]
  },
  "living_off_the_land": {
    "label": "Living Off The Land",
    "mitre_tactic_id": "TA0005",
    "mitre_tactic": "Defense Evasion",
    "mitre_technique_id": "T1218",
    "mitre_technique": "System Binary Proxy Execution",
    "default_severity": "Medium",
    "default_threshold": 1,
    "default_window_minutes": 5,
    "log_source": "Sysmon (Event ID 1), EDR process telemetry",
    "required_fields": ["process_name", "command_line", "parent_process", "host"],
    "false_positives": ["Legitimate IT administration using native utilities", "Software installers invoking regsvr32/rundll32 as designed", "Vendor management tools using bitsadmin/certutil for updates"],
    "coverage_gaps": ["Novel or rare LOLBins not yet catalogued", "Chained LOLBin execution split across multiple low-signal steps", "LOLBin usage on non-Windows platforms (Linux/macOS binaries)"],
    "improvements": ["Correlate LOLBin usage with unusual parent processes", "Add command-line argument anomaly scoring", "Cross-reference with LOLBAS project indicators"]
  },
  "reconnaissance": {
    "label": "Reconnaissance",
    "mitre_tactic_id": "TA0043",
    "mitre_tactic": "Reconnaissance",
    "mitre_technique_id": "T1595",
    "mitre_technique": "Active Scanning",
    "default_severity": "Low",
    "default_threshold": 20,
    "default_window_minutes": 5,
    "log_source": "Firewall/IDS logs, DNS logs, Windows Security (4661, 4662)",
    "required_fields": ["source_ip", "destination_ip", "port", "query"],
    "false_positives": ["Authorized vulnerability scanning by the security team", "Network monitoring tools performing health checks", "Asset discovery/inventory tools"],
    "coverage_gaps": ["Slow, distributed scanning below volumetric thresholds", "Passive reconnaissance with no active network signal", "Reconnaissance performed through legitimate cloud APIs"],
    "improvements": ["Correlate scanning source with subsequent exploitation attempts", "Add allowlisting for known scanner IP ranges", "Track scan pattern diversity (ports/hosts touched) rather than raw volume"]
  }
}
